# deploy/automation-portfolio-api.service.example
[Unit]
Description=Automation Portfolio FastAPI
After=network.target postgresql.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/automation-portfolio/backend
EnvironmentFile=/var/www/automation-portfolio/backend/.env
ExecStart=/var/www/automation-portfolio/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
